# lib/matrix_listener.service
[Unit]
Description=Disinto Matrix Listener
After=network.target dendrite.service

[Service]
Type=simple
ExecStart=/home/admin/disinto/lib/matrix_listener.sh
Restart=always
RestartSec=10
User=admin
WorkingDirectory=/home/admin/disinto

[Install]
WantedBy=multi-user.target
