# lib/matrix_listener.service
[Unit]
Description=Dark Factory Matrix Listener
After=network.target dendrite.service

[Service]
Type=simple
ExecStart=/home/admin/dark-factory/lib/matrix_listener.sh
Restart=always
RestartSec=10
User=admin
WorkingDirectory=/home/admin/dark-factory

[Install]
WantedBy=multi-user.target
